# util/provision/files/server.conf
[program:parms-server]
directory=/vagrant
command=php -S 0.0.0.0:9000 -t www/
autostart=true
autorestart=true
user=vagrant
redirect_stderr=true
stdout_logfile=/vagrant/server.log
